# systemd/xyz.iinuwa.credentialsd.Credentials.service.in
[Unit]
Description=Service for creating and storing user credentials

[Service]
Type=dbus
BusName=xyz.iinuwa.credentialsd.Credentials
ExecStart=@DAEMON_EXECUTABLE@
