# systemd/xyz.iinuwa.credentialsd.UiControl.service.in
[Unit]
Description=Internal helper service for credentialsd

[Service]
Type=dbus
BusName=xyz.iinuwa.credentialsd.UiControl
ExecStart=@UI_EXECUTABLE@
