# systemd/xyz.iinuwa.credentialsd.FlowControl.service.in
[Unit]
Description=Internal helper service for credentialsd

[Service]
Type=dbus
BusName=xyz.iinuwa.credentialsd.FlowControl
ExecStart=@DAEMON_EXECUTABLE@
